# exporter/etc/cron.d/exporter
* 9 * * * root php /code/exporter/main.php
